# backend/statistic/statlang/crontab
0 1 * * * python3 /home/ubuntu/predictTrend/Statistic-Lang/statlang/crawler.py 
